# unisim_cpp.service
[Unit]
Description=start unisim_cpp
[Service]
ExecStart=/usr/local/bin/unisim_server -c can0 
[Install]
WantedBy=multi-user.target
